# root/etc/s6-overlay/s6-rc.d/cdp-proxy/run
#!/command/with-contenv sh
exec socat TCP-LISTEN:9222,fork,reuseaddr TCP:127.0.0.1:9221
