# root/etc/s6-overlay/s6-rc.d/tinyproxy/run
#!/command/with-contenv sh
exec tinyproxy -d -c /app/tinyproxy.conf
