# root/etc/s6-overlay/s6-rc.d/xvnc/run
#!/command/with-contenv sh
exec Xvnc -alwaysshared ${DISPLAY:-:99} -geometry 1920x1080 -depth 24 -rfbport 5900 -SecurityTypes None
